# units/optional/krb5-kadmind.service
[Unit]
Description=Kerberos 5 Administration Server

[Service]
ExecStart=@SBINDIR@/kadmind -nofork

[Install]
WantedBy=multi-user.target
