# units/optional/krb5-kdc.service
[Unit]
Description=Kerberos 5 KDC

[Service]
ExecStart=@SBINDIR@/krb5kdc -n
Restart=always

[Install]
WantedBy=multi-user.target
