# units/optional/krb5-kpropd.service
[Unit]
Description=Kerberos 5 Propagation Server

[Service]
ExecStart=@SBINDIR@/kpropd -S

[Install]
WantedBy=multi-user.target
