# units/optional/sshdat.service
[Unit]
Description=SSH Per-Connection Server

[Service]
ExecStart=-@SBINDIR@/sshd -i
StandardInput=socket
StandardError=journal
