# tests/units/backup.service
[Unit]
Description=System backup service

[Service]
Type=oneshot
ExecStart=/usr/bin/backup-system
