# tests/units/docker.socket
[Unit]
Description=Docker Socket for API

[Socket]
ListenStream=/var/run/docker.sock

[Install]
WantedBy=sockets.target
