# units/reference/console.service
[Unit]
Description=Setup Console Keymap and Font
DefaultDependencies=no
After=sysinit.target
Before=basic.target
ConditionPathExists=@INITD_SYSCONFIG_DIR@/console.conf

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@LIBEXECDIR@/service-scripts/console

[Install]
WantedBy=basic.target
