# units/basic.target
[Unit]
Description=Basic System
Wants=sysinit.target console.service
After=sysinit.target console.service
