# units/default.target
[Unit]
Description=Default Target
Wants=multi-user.target
