# units/echo.service
[Unit]
Description=Echo service (socket activated)

[Service]
Type=simple
ExecStart=/usr/bin/socat STDIN STDOUT
