# units/emergency.target
[Unit]
Description=Emergency Mode
Wants=emergency.service
After=emergency.service
