# units/local-fs.target
[Unit]
Description=Local File Systems
Documentation=man:systemd.special(7)
After=$local_fs
Wants=$local_fs
