# units/multi-user.target
[Unit]
Description=Multi-User System
Requires=basic.target
After=basic.target
Before=graphical.target
OnFailure=basic.target
