# units/multi-user.target
[Unit]
Description=Multi-User System
Wants=basic.target
After=basic.target
Before=graphical.target
