# units/network-online.target
[Unit]
Description=Network is Online
Documentation=man:systemd.special(7)
After=network.target network-online
Wants=network-online
