# units/nfs-server.target
[Unit]
Description=NFS server and services
Documentation=man:systemd.special(7)
After=network.target
