# units/remote-fs.target
[Unit]
Description=Remote File Systems
Documentation=man:systemd.special(7)
After=local-fs.target network.target $remote_fs
Wants=$remote_fs
